# rc.local
!#/bin/bash
/data/openvpn/check_ovpn.sh
